# setup/files/misc/systemd/drbl-clients-nat.service
[Unit]
Description=DRBL NAT service
DefaultDependencies=no
After=network.target

[Service]
Environment=TERM=linux
ExecStart=/etc/init.d/drbl-clients-nat start
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=drbl-clients-nat.service
