# setup/files/misc/systemd/drblthincli.service
[Unit]
Description=DRBL thin client service
DefaultDependencies=no
After=network.target

[Service]
Environment=TERM=linux
ExecStart=/etc/init.d/drblthincli start
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=drbl-clients-nat.service
