# setup/files/DBN/DBN8.0/systemd/system/ocsd-rescue.service
.include /lib/systemd/system/rescue.service

[Unit]
Requires=ocsd-run.service
After=ocsd-run.service

